# init/savvy-report.service
[Unit]
Description=Savvy Report Generator
Documentation=https://github.com/heyajulia/savvy
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=savvy
Group=savvy
ExecStart=/usr/local/bin/savvy report
TimeoutStartSec=120

# Environment file for secrets
EnvironmentFile=/etc/savvy/savvy.env
Environment=STAMP_DIR=/var/lib/savvy/stamps

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RemoveIPC=yes

# Allow writing to stamp directory
ReadWritePaths=/var/lib/savvy/stamps
